# cli/infrastructure/ansible-configurations-bckp/playbooks/roles/node_exporter/templates/com.node-exporter.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.node-exporter</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{ node_exporter_base_dir }}/node_exporter</string>
    <string>--web.listen-address=:{{ node_exporter_port }}</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/node-exporter.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/node-exporter.error.log</string>
</dict>
</plist>
