# ansible-configurations/playbooks/roles/grafana-alloy/templates/grafana-alloy.service.j2
[Unit]
Description=Grafana Alloy - Metrics collection agent
Documentation=https://grafana.com/docs/alloy/
After=network.target

[Service]
User={{ alloy_user }}
Group={{ alloy_user }}
Type=simple
ExecStart={{ alloy_install_dir }}/alloy run \
  --server.http.listen-addr=0.0.0.0:{{ alloy_port }} \
  --storage.path=/var/lib/alloy/data \
  {{ alloy_config_dir }}/config.alloy
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
