# ansible-configurations/playbooks/roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter - Prometheus metrics exporter
Documentation=https://github.com/prometheus/node_exporter
After=network.target

[Service]
User={{ node_exporter_user }}
Group={{ node_exporter_user }}
Type=simple
ExecStart={{ node_exporter_install_dir }}/node_exporter \
  --web.listen-address=:{{ node_exporter_port }}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
