# ansible-configurations/playbooks/roles/openclaw/templates/openclaw.service.j2
[Unit]
Description=OpenClaw AI Assistant
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory=/home/{{ ansible_user }}
ExecStart={{ openclaw_path }}
Restart=on-failure
RestartSec=10s

# Send logs to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
